# rpm/assets/auri.service
[Unit]
Description=auri

[Service]
ExecStart=/usr/bin/auri
User=auri
Group=auri
UMask=007

[Install]
WantedBy=multi-user.target
